# td-stock.service
[Unit]
Description=TD Stock Web Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/td_stock_web
ExecStart=/usr/bin/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target